# debian-amd64-initrd/mkosi.extra.boot/etc/systemd/system/var-lib-extensions.mount
[Mount]
What=/dev/1/extensions
Where=/var/lib/extensions
[Install]
WantedBy=multi-user.target
