# debian-amd64-initrd/mkosi.extra.boot/etc/systemd/system/var-lib-portables.mount
[Mount]
What=/dev/1/portables
Where=/var/lib/portables
[Install]
WantedBy=multi-user.target
